# rootfs/etc/service/nginx/run
#!/bin/sh -e

# pipe stderr to stdout and run nginx
exec 2>&1
exec nginx -g 'daemon off;'